# shell/ipshare.timer
[Unit]
Description=Run IP share periodically

[Timer]
OnUnitActiveSec=5s 

[Install]
WantedBy=timers.target